# console-frontend/deploy/supervisor/supervisor.conf
[supervisord]
nodaemon=true
environment=NODE_ENV="production"

[program:preflight]
command=/bin/sh /opt/app/deploy/scripts/preflight.sh
startsecs=0
autorestart=unexpected

[program:node]
command=node server.js
user=app
autostart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
