# pay-frontend/tools/docker-images/base/fs/etc/cron.d/logrotate
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

*/30 * * * * root /opt/logrotate-rotate.sh
