# pay-frontend/tools/docker-images/base/fs/etc/supervisor/conf.d/syslog.conf
[program:syslog]
command=rsyslogd -n
priority=5
autostart=true
autorestart=true
startsecs=1
startretries=3
stdout_logfile_backups=2
stderr_logfile_backups=2
